# fs-add/etc/inittab
# Boot
::sysinit:/etc/init.d/rcS
# Console
ttyAS0::askfirst:/bin/sh
# Network
::once:/etc/init.d/network
# Time Network Protocol
::once:/etc/init.d/ntpd
# telnetd
::once:/etc/init.d/telnetd
# dropbear
::once:/etc/init.d/dropbear
# oscamd
::once:/etc/init.d/oscamd
# AXE kernel modules and environment
::once:/etc/init.d/axe
# Others
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
